# files/system_manager.init
#!/bin/sh /etc/rc.common

START=95
STOP=10

USE_PROCD=1

start_service() {

    # 네트워크 준비 대기
    while ! ubus call network.interface.lan status | grep -q '"up": true'; do
        sleep 2
    done
    procd_open_instance
    procd_set_param command /usr/bin/system_manager
    # procd_set_param respawn
    procd_close_instance
}
